# .config/systemd/user/obsidian-sync.service
[Unit]
Description=Obsidian Headless Sync
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
ExecStart=/home/gerrit/.local/share/mise/shims/ob sync --path %h/Documents/notes --continuous
Restart=on-failure
RestartSec=10s
[Install]
WantedBy=default.target
